# build/base/usr/local/lib/systemd/system/pureblue.service
[Unit]
Description=pureblue startup script.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pureblue-startup
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
